# cribbage.service
[Unit]
Description=Family Cribbage Tracker
After=network.target

[Service]
User=pi
WorkingDirectory=/home/pi/cribbage-app
Environment="SECRET_KEY=change-this-to-a-random-secret"
Environment="HTTPS=1"
ExecStart=/home/pi/cribbage-app/.venv/bin/gunicorn -w 2 -b 127.0.0.1:5000 app:app
Restart=on-failure

[Install]
WantedBy=multi-user.target
